# deploy/newclaw.service
[Unit]
Description=NewClaw AI Agent Gateway
Documentation=https://github.com/fofo365/newclaw
After=network.target redis.service
Requires=redis.service

[Service]
Type=simple
User=newclaw
Group=newclaw
WorkingDirectory=/opt/newclaw
ExecStart=/opt/newclaw/newclaw gateway
Restart=on-failure
RestartSec=10

# 安全配置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/newclaw /var/log/newclaw /var/run/newclaw

# 资源限制
LimitNOFILE=65535
MemoryMax=1G

# 环境变量
Environment="RUST_LOG=info"
Environment="NEWCLAW_CONFIG=/opt/newclaw/newclaw.toml"

# 标准输出/错误
StandardOutput=journal
StandardError=journal
SyslogIdentifier=newclaw

[Install]
WantedBy=multi-user.target
